# deploy/supervisor/laravel-ami-listener.conf.example
[program:laravel-ami-listener]
process_name=%(program_name)s
directory=/var/www/laravel-crm
command=php artisan ami:listen --reconnect-delay=5
autostart=true
autorestart=true
startsecs=3
startretries=10
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/laravel-crm/storage/logs/ami-listener.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10
stopwaitsecs=15

